# dev/tools/scripts/com.uofa.judge-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Phase 3 Stage 2 daily judge LaunchAgent (spec Part 3.4).
  Install: cp this file to ~/Library/LaunchAgents/com.uofa.judge-daily.plist
           then: launchctl load ~/Library/LaunchAgents/com.uofa.judge-daily.plist
  Fires run_judge_daily.sh at 09:00 local (= 16:00 UTC for PDT), which keeps a
  full pass inside one UTC day so the Gemini RPD cap accounts cleanly. The wrapper
  is idempotent (--resume), so a missed or partial day is recovered the next run.
  RunAtLoad is false: to start Day 1 immediately, run the wrapper once by hand.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.uofa.judge-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/vishnu/Library/CloudStorage/Dropbox/Praxis/uofa_github/dev/tools/scripts/run_judge_daily.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/vishnu/Library/CloudStorage/Dropbox/Praxis/uofa_github/dev/build/adversarial/phase3/production/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vishnu/Library/CloudStorage/Dropbox/Praxis/uofa_github/dev/build/adversarial/phase3/production/launchd.err.log</string>
</dict>
</plist>
